# example/my-app-upstart.conf
description "my-app"

start on (filesystem and net-device-up IFACE!=lo)

chdir /app

console log

respawn
respawn limit 20 30

exec /usr/bin/python3 -m my_app.worker
